# parch-bluetooth/bluetooth.install
post_install() {
	systemctl enable bluetooth
}
