# parch-printer/printer.install
post_install() {
	systemctl enable cups
}
